# compatibility-test.service
[Unit]
Description=.NET Graviton Compatibility Test Application
After=network.target

[Service]
Type=notify
User=ssm-user
WorkingDirectory=/home/ssm-user/Graviton.netApp/src/GravitonBridge.Web
ExecStart=/usr/bin/dotnet run --configuration Release
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=compatibility-test
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target
